# data/gnome-session-failed.service.in
[Unit]
Description=GNOME Session failure lockdown (user)

PartOf=gnome-session-failed.target

# We could do this, but it requires an intermediate target for OnFailure
# handling, so gnome-session-failed checks RUNNING_UNDER_GDM itself
#Conflicts=gnome-session@gnome-login.target
# or in the case of GDM and then not passing --allow-logout
#Requisite=gnome-session@gnome-login.target

[Service]
Type=simple
ExecStart=@libexecdir@/gnome-session-failed --allow-logout
# The fail whale doesn't trigger a shutdown itself, so do it here
ExecStopPost=-@libexecdir@/gnome-session-ctl --shutdown
