# data/gnome-session-shutdown.service.in
[Unit]
Description=Trigger shutdown of a GNOME Session

PartOf=gnome-session.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-session-ctl --shutdown
