# data/gnome-session-post-shutdown.service.in
[Unit]
Description=GNOME Session shutdown cleanup tasks

# Allow exit.target to start even if this unit is started with replace-irreversibly
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=no
ExecStart=@libexecdir@/gnome-session-ctl --restart-dbus
ExecStart=rm %t/gnome-session-leader-fifo
