# data/gnome-session-stable.target
[Unit]
Description=GNOME Session is stable (running for >2 minutes)

Requisite=gnome-session.target
PartOf=gnome-session.target
