# data/gnome-session-signal-init.service.in
[Unit]
Description=Signal initialization done to GNOME Session Manager

PartOf=gnome-session.target

[Service]
Type=oneshot
ExecStart=@libexecdir@/gnome-session-ctl --signal-init
