# data/gnome-session-bus.target
[Unit]
Description=gnome-session is on the bus
StopWhenUnneeded=yes
