# data/gnome-session-manager.target
[Unit]
Description=GNOME Session Manager is ready
RefuseManualStart=true

PartOf=gnome-session-initialized.target
