# data/gnome-session-x11-services.target
[Unit]
Description=GNOME session X11 services
DefaultDependencies=no
